# docker/supervisord/supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock
chmod=0770
chown=nobody:nogroup
# username = mysupervisord
# password = mysupervisordpass

[supervisord]
nodaemon=true
user=root

[program:nginx]
command = nginx -c /etc/nginx/nginx.conf  -g 'daemon off;'
user = root
autostart = true

[program:php-fpm]
command=docker-php-entrypoint php-fpm
user=root
autostart=true
nodaemon=true
autorestart=true

[program:octane]
command=php artisan octane:start --server=swoole --host=0.0.0.0 --port=80 --workers=4 --max-requests=1000
autostart=true
autorestart=true
stdout_logfile=/var/log/octane.log
stderr_logfile=/var/log/octane_error.log